# debian.rabbitmq/3.3.5/files/supervisord.conf
[supervisord]
nodaemon=true

[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server

[program:rabbitmq-clu.radmin]
command=/usr/local/bin/python /rabbitmq-clu.radmin.py
